# s6-overlay/s6-rc.d/koyeb-cli/finish
#!/command/execlineb -S1

if { eltest ${1} -ne 0 -a ${1} -ne 256 }
/run/s6/basedir/bin/halt
